# systemd/.config/systemd/user/wallpaper.service
[Unit]
Description=Randomly change the wallpaper

[Service]
Type=oneshot
ExecStart=wallheaven -c wallpaper

[Install]
WantedBy=timers.target
