# systemd/.config/systemd/user/wallpaper.timer
[Unit]
Description=Change wallpaper every x time

[Timer]
OnActiveSec=0sec
OnUnitActiveSec=10min
